# backup/crontab
15 02 * * * /usr/local/bin/backup-cron.sh >> /var/log/backup.log 2>&1
